# systemd/photo_bot.service
[Unit]
Description=Telegram Photo Frame Bot
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=__USER__
WorkingDirectory=__HOME__
Environment="TELEGRAM_BOT_TOKEN=__BOT_TOKEN__"
ExecStart=/usr/bin/python3 __HOME__/photo_bot.py
Restart=always
RestartSec=10
StandardOutput=append:__HOME__/photo_bot.log
StandardError=append:__HOME__/photo_bot.log
[Install]
WantedBy=multi-user.target
